# ezviz-systemd.net/ezviz-net.service
[Unit]
Description=ezviz.net service

[Service]
Type=notify
ExecStart=/usr/sbin/ezviz-net

[Install]
WantedBy=multi-user.target